# strataapp.service
[Unit]
Description=Strata App Django ASGI Server
After=network.target

[Service]
User=strataapp
Group=strataapp
WorkingDirectory=/opt/strataapp/repo/app/strataapp
EnvironmentFile=/opt/strataapp/strataapp.env

ExecStart=/opt/strataapp/repo/app/.venv/bin/uvicorn strataapp.asgi:application \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 4 \
    --lifespan off

# Process management
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Security
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/strataapp/repo
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
